# build/supervisor/conf.d/app.conf
[program:basil_worker_manager_consume_messages]
command = composer messenger-listen
stopsignal=QUIT
